# HOME-STOW/systemd/.config/systemd/user/i3wm-xsession.target
[Unit]
Description = i3wm Xsession up and running
BindsTo=i3wm-xsession.target
