# HOME-STOW/systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Synchronize imap server

[Service]
ExecStart=/usr/bin/mbsync -a
Type=oneshot

[Install]
WantedBy=default.target
